# files/restic/backup.cron.j2
# /etc/cron.d/backup: crontab entries for backup

8 */6 * * * {{ backup_user }} /home/{{ backup_user }}/backup.bash backup
