# files/restic/backup.cron.j2
# /etc/cron.d/backup: crontab entries for backup

@midnight {{ backup_user }} sudo /home/{{ backup_user }}/backup.bash backup
