# templates/cron.d/backup.j2
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

@daily dime {{ app_root }}/app backup
